# deploy/launchd/com.clawxrouter.privacy-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clawxrouter.privacy-proxy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/share/desensitize-proxy/DesensitizeProxy.AspNetCore</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/usr/local/share/desensitize-proxy</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ASPNETCORE_ENVIRONMENT</key>
    <string>Production</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/desensitize-proxy.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/desensitize-proxy.err.log</string>
</dict>
</plist>
